# services/classify-service.service
[Unit]
Description=Deep Academic Research Classification Service
After=network.target mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/deep-academic-research
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/path/to/deep-academic-research"
ExecStart=/usr/bin/python3 /path/to/deep-academic-research/services/classify_service.py --config /path/to/deep-academic-research/config.rbase.yaml --workers 4
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=classify-service

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 安全选项
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

