# scripts/rbase-mark-api.service
[Unit]
Description=Rbase API Server
After=network.target
Wants=network.target

[Service]
Type=exec
User=rbase
Group=rbase
WorkingDirectory=/opt/rbase/rbase-mark
Environment=PYTHONPATH=/opt/rbase/drbase-mark
ExecStart=/opt/rbase/drbase-mark/venv/bin/python /opt/rbase/drbase-mark/scripts/start_api_server.py --workers 4 --loop asyncio
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rbase-api

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/rbase/drbase-mark/logs /opt/rbase/drbase-mark/data

[Install]
WantedBy=multi-user.target 